# mm-server.service
[Unit]
Description=MagicMirror NodeJs Backend

[Service]
WorkingDirectory=%h/MagicMirror
ExecStart=node serveronly

